# services/wifi-hotspot/wifi-hotspot.service
[Unit]
Description=Cinemate Wi-Fi Hotspot Service
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /usr/local/bin/wifi-hotspot.py
Restart=always
RestartSec=5
User=root

[Install]
WantedBy=multi-user.target
