# services/cinemate-autostart/cinemate-autostart.service
[Unit]
Description=Cinemate Application Service
After=network.target plymouth-quit-wait.service
Before=getty@tty1.service
Conflicts=getty@tty1.service

[Service]
Type=notify
NotifyAccess=main
Environment=CINEMATE_STARTUP_FAILURE_FILE=/home/pi/.cache/cinemate/startup-failure.ansi
StandardInput=tty-force
TTYPath=/dev/tty1
TTYReset=yes
TTYVHangup=yes
TTYVTDisallocate=no
ExecStartPre=/usr/local/bin/camera-ready.sh
ExecStartPre=/bin/sleep 5
ExecStart=/home/pi/.cinemate-env/bin/python3 /home/pi/cinemate/src/main.py
ExecStopPost=-+/usr/local/bin/cinemate-startup-failure-display.sh
ExecStopPost=-+/bin/systemctl start getty@tty1.service
WorkingDirectory=/home/pi/cinemate
User=pi
Group=pi
SyslogIdentifier=cinemate-autostart
LogRateLimitIntervalSec=30s
LogRateLimitBurst=20
TimeoutStopSec=5
KillSignal=SIGINT

[Install]
WantedBy=multi-user.target
